# firefox-ram.service
[Unit]
Description=Firefox profile memory cache

[Install]
WantedBy=default.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/ram-firefox %i
ExecStop=/usr/local/bin/ram-firefox %i
